# files/system/etc/profile.d/open.sh
#!/usr/bin/env -S sh
alias open="xdg-open &>/dev/null"