# system_files/etc/profile.d/fastfetch.sh
alias fastfetch='fastfetch --logo "$(find /usr/share/ublue-os/bluefin-logos/symbols -type f | shuf -n1)" --color "$(ublue-bling-fastfetch)" --config /etc/ublue-os/fastfetch.json'
alias neofetch='fastfetch'
alias neowofetch='fastfetch'