# system_files/etc/profile.d/open.sh
#!/usr/bin/env -S bash
alias open='xdg-open &>/dev/null'