# system_files/etc/profile.d/ublue-motd.sh
#!/usr/bin/env -S bash

[[ ! -e "$HOME"/.config/no-show-user-motd ]] && /usr/libexec/ublue-motd